# demo/images/src/porthos/supervisord.conf
[supervisord]
nodaemon=true

[program:ioam-register]
command=/apps/ioam_register
autorestart=false
startretries=0
