# demo/images/src/porthos/supervisord.conf
[supervisord]
nodaemon=true

[program:empty]
command=/apps/empty
autorestart=false
startretries=0
