# demo/docker/porthos/supervisord.conf
[supervisord]
nodaemon=true

[program:ioam-config]
startsecs=0
command=/sbin/ip ioam namespace add 123

[program:app]
command=tail -f /dev/null
